# deployment/celestia-backend.service
[Unit]
Description=Celestia Ground Station Backend
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=CE
Group=CE
WorkingDirectory=/home/CE/celestia-groundstation
Environment="PATH=/home/CE/celestia-groundstation/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/CE/celestia-groundstation/venv/bin/python -m backend.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celestia-backend

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
